# openwrt/data/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
mwan3 restart exit 0
smartdns restart exit 0
